# supervisord.conf
;supervisor config
# [unix_http_server]
# file = /var/run/supervisord.pid
# chown = www-data

[supervisord]
logfile=/proc/1/fd/1
logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
stdout_logfile=/proc/1/fd/1
stderr_logfile=/proc/1/fd/1
stdout_maxbytes=0
stderr_maxbytes=0
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
# nodaemon=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stderr_logfile=/proc/1/fd/1
stopwaitsecs=3600
environment=DB_HOST="%(ENV_DB_HOST)s",DB_DATABASE="%(ENV_DB_DATABASE)s",DB_USERNAME="%(ENV_DB_USERNAME)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",REDIS_HOST="redis"

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command= /bin/sh -c "while [ true ]; do (php /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stderr_logfile=/proc/1/fd/1
stopwaitsecs=3600
environment=DB_HOST="%(ENV_DB_HOST)s",DB_DATABASE="%(ENV_DB_DATABASE)s",DB_USERNAME="%(ENV_DB_USERNAME)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",REDIS_HOST="redis"
